# src/utils/crontab/dev-crontab
# min hour day month weekday command
* * * * * ts-node -e 'require("/home/app/src/subscribers/checkFTP.ts").check("/home/app/");' >> /data/logs/cron.log 2>&1
